# systemd/peering-manager_peeringdb-sync.timer
[Unit]
Description=Sync PeeringDB at 2:30

[Timer]
OnCalendar=02:30:00

[Install]
WantedBy=timers.target
