# systemd/peering-manager_configuration-deployment.timer
[Unit]
Description=Push router configuration every hour 5 minutes before full hour

[Timer]
OnCalendar=*:55:00

[Install]
WantedBy=timers.target
