# systemd/peering-manager_prefix-fetch.timer
[Unit]
Description=Fetch IRR AS-SET prefixes at 4:30

[Timer]
OnCalendar=04:30:00

[Install]
WantedBy=timers.target
